# docker/app/crontab
1 * * * * /usr/bin/python3 /opt/papermerge/manage.py txt2db  >> /var/log/cron.log 2>&1
1 * * * * /usr/bin/python3  /opt/papermerge/manage.py update_fts  >> /var/log/cron.log 2>&1
