# config/templates/etc/systemd/papermerge.target
{{notice}}

[Unit]
Description=Papermerge Targets
Requires=papermerge.service \
     mg_nginx.service \
     worker.service \
