# config/templates/etc/systemd/papermerge.target
[Unit]
Description=Papermerge Targets
Requires=papermerge.service pm_nginx.service celery.service