# config/templates/etc/systemd/papermerge.target
{{notice}}

[Unit]
Description=Papermerge Targets
Requires=papermerge.service pm_nginx.service pmworker.service