# config/templates/etc/systemd/papermerge.target
{{notice}}

[Unit]
Description=Papermerge Targets
Requires=papermerge.service \
     pm_nginx \
     pmworker \
     txt2db \
     update_fts