# docker/1.2/app/crontab
* * * * * /opt/txt2db.sh
* * * * * /opt/update_fts.sh
* * * * * /opt/hello.sh