# docker/etc/txt2db.service
[Unit]
Description=txt2db service

[Service]
Type=oneshot
WorkingDirectory=/opt/papermerge/
ExecStart=python /opt/papermerge/manage.py txt2db 